# systemd/pkg-listn.path
[Unit]
Description=watch ~/.config/pkg-listn/packages
PartOf=graphical-session.target

[Path]
PathModified=%h/.config/pkg-listn/packages

[Install]
WantedBy=graphical-session.target

