# systemd/pkg-listn.service
[Unit]
Description=install or remove packages interactively

[Service]
Type=oneshot
ExecStart=pkg-listn
